# _archive/legacy-dirs/deploy/systemd/sultan-node.service
[Unit]
Description=Sultan L1 Blockchain Node
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=sultan
Group=sultan
WorkingDirectory=/opt/sultan
Environment="RUST_LOG=info"

# Node configuration
ExecStart=/opt/sultan/bin/sultan-node \
  --enable-sharding \
  --shard-count 100 \
  --tx-per-shard 10000 \
  --block-time 5 \
  --data-dir /var/lib/sultan/data \
  --rpc-addr 127.0.0.1:26657 \
  --validator \
  --validator-address sultan-validator-1 \
  --validator-stake 1000000

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/sultan

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Restart policy
Restart=always
RestartSec=10
StartLimitInterval=200
StartLimitBurst=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sultan-node

[Install]
WantedBy=multi-user.target
